# deploy/wol.service
# DEPRECATED: wol.service
# This file is kept for historical/reference purposes only.
# DO NOT enable or use this unit. The canonical unit to use is `wakeonlan.service`.
#
# To remove the old unit from a system where it was previously deployed, run:
#   sudo systemctl disable --now wol.service || true
#   sudo rm -f /etc/systemd/system/wol.service || true
#   sudo systemctl daemon-reload
#
# If you need the recommended unit file, use the provided `deploy/wakeonlan.service`.

# Example systemd service for Raspberry Pi (adjust paths and user)
[Unit]
Description=WakeOnLan Gunicorn (example service)
After=network.target

[Service]
# Run as the 'wol' user (create user 'wol' and place code in /home/wol/Wake-on-lan)
User=wol
Group=wol
WorkingDirectory=/home/wol/Wake-on-lan
# If you use a virtualenv, point ExecStart to the venv gunicorn binary below.
# Ensure the PATH contains the venv bin so subprocesses can find python if needed.
Environment=PATH=/home/wol/Wake-on-lan/.venv/bin:/usr/local/bin:/usr/bin
# Allow binding to port 80 without running as root
AmbientCapabilities=CAP_NET_BIND_SERVICE

# Gunicorn exec: adjust the path to the gunicorn binary installed in the venv
ExecStart=/home/wol/Wake-on-lan/.venv/bin/gunicorn \
  --chdir /home/wol/Wake-on-lan \
  --worker-class gthread \
  --workers 2 \
  --threads 4 \
  --timeout 120 \
  --graceful-timeout 30 \
  --bind 0.0.0.0:80 \
  --access-logfile - \
  --error-logfile - \
  wol_app:app

# Restart and basic hardening
Restart=on-failure
RestartSec=5s
# Give workers a bit more time to stop gracefully
TimeoutStopSec=60s
KillMode=mixed
PrivateTmp=yes
ProtectSystem=full
ProtectHome=yes
NoNewPrivileges=yes
PrivateDevices=yes
LimitNOFILE=4096

# Let systemd capture stdout/stderr into the journal
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
